# config/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/home/admin/wagtail-italia.it/run/gunicorn.sock

[Install]
WantedBy=sockets.target